# deploy/systemd/exchangekit-frontend.service
[Unit]
Description=ExchangeKit frontend (Next.js)
After=network.target exchangekit-backend.service
Wants=exchangekit-backend.service

[Service]
Type=exec
User=exchangekit
Group=exchangekit
WorkingDirectory=__REPO_DIR__/frontend
EnvironmentFile=__REPO_DIR__/.env
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/node __REPO_DIR__/frontend/node_modules/next/dist/bin/next start -p 3000
Restart=always
RestartSec=3
NoNewPrivileges=true
ProtectSystem=full
PrivateTmp=true

[Install]
WantedBy=multi-user.target
